# host/ubuntu2404/cron.service
[Unit]
Description=Staticpatch Cron

[Service]
ExecStart=/bin/bash -c "/home/staticpatch/run_manage.sh cron"
User=staticpatch
Group=staticpatch
